# provisioning/files/ngrok-init.conf
# Startup the ngrok proxy

description "ngrok proxy instance"

start on runlevel [2345]
stop on runlevel [!2345]

exec /home/skipio/ngrok http 3000